# qmd/luan.qmd-reindex.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>luan.qmd-reindex</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/luan/.local/bin/qmd-reindex</string>
    </array>

    <key>StandardErrorPath</key>
    <string>/tmp/qmd-reindex.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/qmd-reindex.out</string>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
